# deploy/asi-hub.service
[Unit]
Description=ASI Deploy Hub API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/asi-deploy-hub
Environment=ASI_DEPLOY_DB=DRIVER={ODBC Driver 18 for SQL Server};SERVER=localhost,1433;DATABASE=ASIDeployHub;UID=sa;PWD=changeme;Encrypt=no;TrustServerCertificate=yes;
Environment=ASI_ARTIFACT_STORE=/opt/data/asi-artifacts
ExecStart=/usr/bin/python3 /opt/asi-deploy-hub/src/hub.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
